# dist/src/main/unix/scripts/postinst
#!/bin/bash

set -e

GROUP='stratio'
USER='sparta'

case "$1" in
    configure)
    if ! getent group ${GROUP} >/dev/null 2>&1; then
       groupadd ${GROUP}
    fi
    if ! getent passwd ${USER} >/dev/null 2>&1; then
       useradd -M -d /opt/sds/sparta -g ${GROUP} -s /sbin/nologin ${USER}
    fi

    chown root:root /etc/init.d/sparta
    chown root:root /etc/default/sparta-functions
    chown root:root /etc/default/sparta-variables

    chown -R ${USER}:${GROUP} /etc/sds/sparta
    chown -R ${USER}:${GROUP} /var/log/sds/sparta
    chown -R ${USER}:${GROUP} /opt/sds/sparta

    chmod 755 /etc/init.d/sparta
    chmod 755 /etc/sds/sparta
    chmod 755 /var/log/sds/sparta
    chmod 755 /opt/sds/sparta
    chmod 750 /opt/sds/sparta/bin/*
    
    chmod 644 /etc/default/sparta-functions
    chmod 644 /etc/default/sparta-variables

    update-rc.d sparta defaults
  ;;
esac

exit 0
