# dist/src/main/unix/scripts/postinst
#!/bin/sh

USER=stratio
GROUP=stratio
DIR_INSTALL=/opt/sds/sparkta
DIR_PID=/var/run/sds
DIR_LOG=/var/log/sds/sparkta
PID_FILE=$PID_DIR/sparkta.pid
CONF_FILE_BACKUP=$DIR_INSTALL/conf.tar.gz
LOG_FILE=$DIR_LOG/sparkta.out
LOG_ERR_FILE=$DIR_LOG/sparkta.err

###
# UTILs
###
restoreAndPurge () {
    if [ -f "$1" ]; then
        echo "Restoring $1"
        tar xzf $1 > /dev/null 2>&1
        echo "Removing $1"
        rm $1
    fi
}

###
# EXECUTION CONFIGURATION
###
/bin/touch $PID_FILE
/bin/touch $LOG_FILE
/bin/touch $LOG_ERR_FILE
chown -R ${USER}:${GROUP} $DIR_INSTALL
chown -R ${USER}:${GROUP} $DIR_PID
chown -R ${USER}:${GROUP} $DIR_LOG
/bin/chmod 775 $DIR_INSTALL
/bin/chmod 755 $DIR_PID
/bin/chmod 755 $DIR_LOG
/bin/chmod 777 $PID_FILE
/bin/chmod 755 $LOG_FILE
/bin/chmod 755 $LOG_ERR_FILE
/bin/chmod 755 $DIR_INSTALL/bin/run
/bin/chmod 755 /etc/init.d/sparkta
echo "Adding sparkta service"
update-rc.d sparkta defaults > /dev/null 2>&1

###
# RESTORE conf
###
restoreAndPurge $CONF_FILE_BACKUP

exit 0