# dist/src/main/unix/scripts/postrm
#!/bin/sh

DIR_PID=/var/run/sds

update-rc.d -f sparkta remove > /dev/null 2>&1

exit 0