# dist/src/main/unix/scripts/postinst
#!/bin/bash

set -e

GROUP='stratio'
USER='sparkta'

case "$1" in
    configure)
    if ! getent group ${GROUP} >/dev/null 2>&1; then
       groupadd ${GROUP}
    fi
    if ! getent passwd ${USER} >/dev/null 2>&1; then
       useradd -M -d /opt/sds/sparkta -g ${GROUP} -s /sbin/nologin ${USER}
    fi

    chown root:root /etc/init.d/sparkta

    chown -R ${USER}:${GROUP} /etc/sds/sparkta

    chown ${USER}:${GROUP} /etc/default/sparkta-functions
    chown ${USER}:${GROUP} /etc/default/sparkta-variables

    chown -R ${USER}:${GROUP} /var/log/sds/sparkta
    chown -R ${USER}:${GROUP} /opt/sds/sparkta

    chmod 755 /etc/init.d/sparkta

    chmod 750 /etc/sds/sparkta
    chmod 775 /var/log/sds/sparkta
    chmod 750 /opt/sds/sparkta
    chmod 750 /opt/sds/sparkta/bin/*

    update-rc.d sparkta defaults
  ;;
esac

exit 0