# dist/src/main/unix/scripts/prerm
#!/bin/sh
 
set -e
 
/etc/init.d/sparkta stop
if [ "$1" = remove ]; then
    invoke-rc.d sparkta stop || true
fi