# dist/src/main/unix/scripts/prerm
#!/bin/sh
 
set -e
 
/etc/init.d/sparta stop
if [ "$1" = remove ]; then
    invoke-rc.d sparta stop || true
fi